# packages/portinus/portinus-1.0.17.tar.gz/portinus-1.0.17/portinus/templates/instance.service
[Unit]
Description={{ name }} service
After=network.target docker.service
PartOf=docker.service

[Service]
Type=simple
ExecStart={{ start_command }}
ExecStop={{ stop_command }}

[Install]
WantedBy=multi-user.target
