# packages/portinus/portinus-1.0.17.tar.gz/portinus-1.0.17/portinus/templates/monitor.service
[Unit]
Description=Monitor module for {{ name }} service
After=network.target docker.service

[Service]
Type=oneshot
ExecStart={{ portinus_monitor_path }} check --name {{ name }}

