# packages/portinus/portinus-1.0.17.tar.gz/portinus-1.0.17/portinus/templates/restart.timer
[Unit]
Description=Restart timer for {{ name }} service

[Timer]
OnCalendar={{ restart_schedule }}
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
