# packages/portinus/portinus-1.0.17.tar.gz/portinus-1.0.17/portinus/templates/restart.service
[Unit]
Description=Restart module for {{ name }} service
After=network.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart {{ service_name }}
